# systemd/scout-audit.service
[Unit]
Description=Wildcatter Scorer Audit — Monthly drift check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/opportunity-scout
EnvironmentFile=/opt/opportunity-scout/.env
ExecStart=/bin/bash -c 'cd /opt/opportunity-scout && /opt/opportunity-scout/venv/bin/python -m src.cli scorer_audit'
TimeoutStartSec=600
MemoryMax=256M
Nice=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scout-audit

[Install]
WantedBy=multi-user.target
