# systemd/scout-daily.timer
[Unit]
Description=OpportunityScout Daily Scan Timer

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
